# tjcsl-discord-bot.service
[Unit]
Description=TJCSL Discord Bot
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/krishnan/tjcsl-discord-bot
ExecStartPre=/bin/sh -c 'until host discord.com; do sleep 1; done'
ExecStart=python3.11 -m poetry run python -u /home/krishnan/tjcsl-discord-bot/main.py
Restart=on-abnormal  # Restart on ExecStartPre timeout
RestartSec=60  # Wait a minute to restart
User=krishnan
Group=krishnan

[Install]
WantedBy=multi-user.target
